# luci-theme-leo/root/etc/uci-defaults/30_luci-theme-leo
#!/bin/sh
uci batch <<-EOF
	set luci.themes.AdvancedTomatoleo=/luci-static/leo
	commit luci
EOF
exit 0
